# cosmetic/dotfiles/config/autostart/btop-fixed.desktop
[Desktop Entry]
Type=Application
Exec=alacritty --class "btop-monitor" -e btop
Name=Btop Fixed
X-KDE-StartupFeedback=false
